# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor_logs/supervisord.log
user=root

[program:redis]
command=redis-server

[program:dp]
command=/dproxy/start.sh
